# systemd/pg-backup.service
[Unit]
Description=PostgreSQL Backup — all user databases
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=labmon
Group=labmon
ExecStart=/root/LabDoctorM/projects/lab-monitoring/scripts/pg_backup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
